# ecg-web.service
[Unit]
Description=ECG DICOM Converter API
After=network.target

[Service]
Type=exec
User=ecg-web
Group=ecg-web
WorkingDirectory=/opt/ecg-web
EnvironmentFile=/opt/ecg-web/.env
Environment=MPLCONFIGDIR=/var/cache/ecg-web/matplotlib
ExecStart=/opt/ecg-web/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001 --proxy-headers
Restart=on-failure
RestartSec=5

# Hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/ecg-web /var/cache/ecg-web

[Install]
WantedBy=multi-user.target
